# systemd/borg-backup.timer
[Unit]
Description=continuous invocation of borg backup

[Timer]
OnActiveSec=30s
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
